# restic-server/systemd/backup-check.timer
[Unit]
Description=check restic repositories

[Timer]
OnCalendar=05:00
Unit=backup-check.target

[Install]
WantedBy=timers.target
