# restic-server/systemd/backup-prune.target
[Unit]
Description=Remove old backups
StopWhenUnneeded=yes
Wants=backup-prune@alpha.service
Wants=backup-prune@zeta.service
