# restic-server/systemd/backup-cloud-prune.target
[Unit]
Description=Remove old cloud backups
Wants=backup-cloud-prune@alpha.service
Wants=backup-cloud-prune@zeta.service
