# restic-server/systemd/backup-prune.timer
[Unit]
Description=Remove old backups

[Timer]
OnCalendar=Sat 05:15
Unit=backup-prune.target

[Install]
WantedBy=timers.target
