# restic-server/systemd/backup-smb.timer
[Unit]
Description=Periodically backup samba

[Timer]
OnCalendar=02:15

[Install]
WantedBy=timers.target
