# restic-server/systemd/backup-cloud-copy.target
[Unit]
Description=Copy backups to the cloud
Wants=backup-cloud-copy@alpha.service
Wants=backup-cloud-copy@zeta.service
