# restic-server/systemd/backup-check.target
[Unit]
Description=check restic repositories
StopWhenUnneeded=yes
Wants=backup-check@alpha.service
Wants=backup-check@zeta.service
