# restic-server/systemd/backup-cloud-copy.timer
[Unit]
Description=Copy backups to the cloud

[Timer]
OnCalendar=05:25
Unit=backup-cloud-copy.target

[Install]
WantedBy=timers.target
